# resources/docker/cron/image/crontab
* * * * * /usr/local/bin/php /app/bin/yii hello
